# .config/systemd/user/sway-session.target.wants/kanshi.service
[Unit]
Description=mako notofication daemon
PartOf=sway-session.target
After=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kanshi
Restart=on-failure

[Install]
WantedBy=sway-session.target
